# deploy/block-backend.service
[Unit]
Description=Block — backend de voto electronico (NestJS)
# El backend gestiona la red Fabric con comandos docker: no puede arrancar
# antes que el demonio.
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=__USER__
WorkingDirectory=__APP_DIR__/backend
EnvironmentFile=__APP_DIR__/backend/.env
ExecStart=/usr/bin/node dist/main.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
